# systemd/vortexl2-forward-daemon.service
[Unit]
Description=VortexL2 Port Forward Daemon (nftables)
After=network.target vortexl2-tunnel.service nftables.service
Requires=network.target
Wants=nftables.service
Documentation=file:///usr/share/doc/vortexl2/README.md

[Service]
Type=simple
WorkingDirectory=/opt/vortexl2
Environment=PYTHONPATH=/opt/vortexl2
ExecStart=/usr/bin/python3 -m vortexl2.forward_daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vortexl2-forward

# Security
NoNewPrivileges=false
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
